# deploy/cc-bridge.service
[Unit]
Description=cc-bridge — persistent Claude Code sessions over OpenAI-compatible HTTP
After=network-online.target
Wants=network-online.target

[Service]
Type=simple

# Replace @CC_BRIDGE_DIR@ with the absolute path to your cc-bridge clone:
#   sed "s|@CC_BRIDGE_DIR@|$PWD|g" deploy/cc-bridge.service \
#     > ~/.config/systemd/user/cc-bridge.service
#   systemctl --user daemon-reload && systemctl --user enable --now cc-bridge
WorkingDirectory=@CC_BRIDGE_DIR@
ExecStart=/usr/bin/env node @CC_BRIDGE_DIR@/src/daemon.js
Restart=on-failure
RestartSec=5
TimeoutStopSec=30
TimeoutStartSec=30
SuccessExitStatus=0 143
KillMode=control-group

# Logging — journal, tagged for easy filtering.
StandardOutput=journal
StandardError=journal
SyslogIdentifier=cc-bridge

# Secrets + runtime config come from the bridge .env. Keep this unit
# secret-free — everything tunable lives in `.env` (see .env.example).
EnvironmentFile=@CC_BRIDGE_DIR@/.env

# Constants that must not be overridable from .env.
Environment=NODE_ENV=production

# Claude Code uses OAuth from ~/.claude/.credentials.json (Max sub, $0).
# ANTHROPIC_API_KEY MUST be empty — it would override OAuth and start
# billing the pay-per-token endpoint.
Environment=ANTHROPIC_API_KEY=

# Intentionally NO sd_notify / WatchdogSec. Liveness is enforced by
# Restart=on-failure plus the daemon's per-child supervisor (see
# src/supervisor.js).

[Install]
WantedBy=default.target
